# systemd/entilldisplay.service
[Unit]
Description=entilldisplay — self-hostad signage-player
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=eriks
SupplementaryGroups=video render
Environment=POLL=60
# Supervisorn kör player:n + hanterar OTA-uppdatering och rollback.
ExecStart=/opt/entilldisplay/bin/supervisor.sh vagg5
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
